# ansible/templates/supervisor_swe30004.conf
[program:swe30004]
command=/app/.virtualenv/bin/python /app/main.py
directory=/app
user=app
autostart=true
autorestart=true
stderr_logfile=/var/log/swe30004.err.log
stdout_logfile=/var/log/swe30004.out.log
environment=APP_SETTINGS=config.ProductionConfig
